# crontab
0 0 * * * python /schooljobsvicscraper.py >> /logfile
