# crontab
0 0 * * * python /rolscraper.py >> /logfile
